# deploy/systemd/industrial-automation-quoter-watchdog.service
[Unit]
Description=Industrial Automation Quoter Health Watchdog
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/saulordz/projects
Environment=HOME=/home/saulordz
Environment=PATH=/home/saulordz/n/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/saulordz/projects/scripts/health-watchdog.sh

